# production/build/podman/akeneo/ak5/files/queue/pim_job_queue.service
[Unit]
Description=Akeneo PIM Job Queue Service (#%i)
After=apache2.service

[Service]
Type=simple
User=build
WorkingDirectory=/build/app
ExecStart=/build/app/bin/console akeneo:batch:job-queue-consumer-daemon --env=prod
Restart=always

[Install]
WantedBy=multi-user.target